# scripts/crontab
# Crontab para Liga de Lincoln
# Scripts en /home/gallardo/Documentos/ligadelincoln/scripts/

# Scraper de resultados - Sáb/Dom cada 15min, 14-21hs Argentina (UTC 17-23 + 0)
*/15 17-23 * * 6,0 /home/gallardo/Documentos/ligadelincoln/scripts/run_scraper_resultados.sh >> /home/gallardo/logs/scraper_resultados.log 2>&1
*/15 0 * * 6,0 /home/gallardo/Documentos/ligadelincoln/scripts/run_scraper_resultados.sh >> /home/gallardo/logs/scraper_resultados.log 2>&1

# Placas de resultados - Domingos 22:15 Argentina (UTC 01:15)
15 1 * * 0 /home/gallardo/Documentos/ligadelincoln/scripts/generar_placas_resultados.sh >> /home/gallardo/logs/generador_placas.log 2>&1
